# hass-stack/config/mariadb-init/crontab
# MariaDB cron jobs
0 1 * * * /docker-entrypoint-initdb.d/backup.sh